# provisioning/templates/supervisor.conf
[program:{{ project_name }}]
command={{ project_bin }}/runserver
user={{ project_user }}
environment=PATH="{{ project_bin }}"
